# src/standalone/contents/etc/service/dsdnsd/run
#!/bin/bash
exec /opt/dnssafety/bin/dsdnsd
